# PieBrain/ChedWeb/frontend/cheddar-frontend.service
[Unit]
Description=Cheddar Frontend (React Production Build)
After=network-online.target cheddar-backend.service
Wants=network-online.target
Requires=cheddar-backend.service

[Service]
Type=simple
User=adamprobert
WorkingDirectory=/home/adamprobert/Cheddar/PieBrain/ChedWeb/frontend
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Extended timeout for npm install and build on Raspberry Pi
TimeoutStartSec=300

# Pre-start: Git sync from main branch
ExecStartPre=/bin/bash -c 'cd /home/adamprobert/Cheddar && git fetch origin && git checkout main && git reset --hard origin/main'

# Install dependencies and build
ExecStartPre=/bin/bash -c 'cd /home/adamprobert/Cheddar/PieBrain/ChedWeb/frontend && npm install'
ExecStartPre=/bin/bash -c 'cd /home/adamprobert/Cheddar/PieBrain/ChedWeb/frontend && npm run build'

# Serve production build on port 3000
ExecStart=/usr/bin/npx serve -s dist -l 3000

# Restart configuration: 5 attempts within 10 minutes
Restart=on-failure
RestartSec=10
StartLimitInterval=600
StartLimitBurst=5

# Logging (journalctl -u cheddar-frontend -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-frontend

[Install]
WantedBy=multi-user.target
