# PieBrain/ChedWeb/frontend/cheddar-frontend.service
[Unit]
Description=Cheddar Frontend (React Production Build)
After=network-online.target cheddar-backend.service
Wants=network-online.target
# Wants, not Requires: this only serves static files and the browser talks to
# the backend on :8000 directly. Under Requires, a failed backend start job
# cancelled this one outright and it stayed dead even after the backend's
# Restart= brought it back, taking the UI down with it.
Wants=cheddar-backend.service

[Service]
Type=simple
User=adamprobert
WorkingDirectory=/home/adamprobert/Cheddar/PieBrain/ChedWeb/frontend
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# The frontend is intentionally NOT built on boot. Running
# `npm install && npm run build` on a Raspberry Pi 3B takes longer than any
# reasonable start timeout, which caused a start-pre timeout -> restart loop
# and the UI never came up on :3000. The production bundle in dist/ is built
# ahead of time instead:
#   - once by install_services.sh during setup, and
#   - after any frontend change via ./build-frontend.sh
# This service only serves the pre-built dist/.
# (npx may download `serve` on first run, so allow a little startup headroom.)
TimeoutStartSec=120

# Serve the pre-built production bundle on port 3000
ExecStart=/usr/bin/npx serve -s dist -l 3000

# Restart configuration:
Restart=on-failure
RestartSec=10

# Logging (journalctl -u cheddar-frontend -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cheddar-frontend

[Install]
WantedBy=multi-user.target
